# .claude/launchd-templates/com.aion.jarvis-cost-watcher.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.aion.jarvis-cost-watcher</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>/Users/nathanielcannon/Claude/Jarvis/.claude/scripts/cost-anomaly-watcher.sh</string>
    </array>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PROJECT_DIR</key>
        <string>/Users/nathanielcannon/Claude/Jarvis</string>
        <key>PATH</key>
        <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
        <key>HOME</key>
        <string>/Users/nathanielcannon</string>
        <key>COST_INTERVAL</key>
        <string>300</string>
    </dict>

    <key>WorkingDirectory</key>
    <string>/Users/nathanielcannon/Claude/Jarvis</string>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <true/>

    <key>StandardOutPath</key>
    <string>/Users/nathanielcannon/Claude/Jarvis/.claude/logs/cost-anomaly-watcher.stdout.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/nathanielcannon/Claude/Jarvis/.claude/logs/cost-anomaly-watcher.stderr.log</string>

    <key>LimitLoadToSessionType</key>
    <string>Aqua</string>
</dict>
</plist>
